# config/supervisor-pure-ftpd.conf
[program:pure-ftpd]
command=/usr/sbin/pure-ftpd-postgresql-virtualchroot -A -l pgsql:/etc/pure-ftpd/db/postgresql.conf -j -p 30000:30009 -P EXTERNAL_IP
