# config/supervisor-pure-ftpd.conf
[program:pure-ftpd]
command=/usr/sbin/pure-ftpd-mysql-virtualchroot -A -l mysql:/etc/pure-ftpd/db/mysql.conf -j -p 30000:30009 -P EXTERNAL_IP
